# script/vllm_ocr_worker.service
[Unit]
Description=Celery vLLM OCR Worker Service
After=network.target rabbitmq-server.service
Wants=network.target rabbitmq-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/Pre-trained

# Set environment variables
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/admin/Pre-trained/.env
Environment="TRITON_PTXAS_PATH=/usr/local/cuda/bin/ptxas"

ExecStart=/home/admin/Pre-trained/.venv/bin/python3.12 -m celery -A src.gpu_server.celery_app.celery_app worker -n vllm_ocr_worker --pool=solo -Q vllm_ocr --loglevel=info

Restart=always
RestartSec=5
KillMode=mixed

StandardOutput=append:/home/admin/Pre-trained/logs/vllm_ocr_worker.log
StandardError=append:/home/admin/Pre-trained/logs/vllm_ocr_worker.log

[Install]
WantedBy=multi-user.target

